# deploy/systemd/portfolio-activity-sync.service
[Unit]
Description=Refresh the server-side portfolio activity cache
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/portfolio-activity-sync
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=CI=1
Environment=NO_COLOR=1
Environment=ACTIVITY_CACHE_PATH=%h/.local/share/rag-playground/activity.json
UMask=0077
ExecStartPre=/usr/bin/env codex doctor --json --no-color
ExecStart=/usr/bin/python3 %h/portfolio-activity-sync/scripts/sync_portfolio_activity.py
TimeoutStartSec=15min
NoNewPrivileges=true
PrivateTmp=true
